# src/create_fastapi/templates/deploy/supervisor/{{ project_name }}.conf
; supervisor 进程配置（部署前请按实际环境修改 directory / user / 路径）
[program:{{ project_name }}]
directory=/srv/{{ project_name }}
command=/srv/{{ project_name }}/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/{{ project_name }}/logs/uvicorn.log
environment=APP_ENV="production"
{% if use_celery %}
[program:{{ project_name }}-worker]
directory=/srv/{{ project_name }}
command=/srv/{{ project_name }}/.venv/bin/celery -A app.celery_app:celery_app worker -l info
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/{{ project_name }}/logs/celery.log
environment=APP_ENV="production"
{% endif %}
